# config/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

exec su - <%= user %> -c 'cd <%= engine.root %>; export APP_ENV=production; <%= process.command %> >> <%= log %>/<%=name%>-<%=num%>.log 2>&1'
